# supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700
chown=chromeuser:chromeuser

[supervisord]
nodaemon=true
user=root
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -ac +extension GLX +render -noreset
user=chromeuser
autorestart=true
priority=10
environment=HOME="/home/chromeuser",USER="chromeuser",LOGNAME="chromeuser"

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -forever -shared -listen 0.0.0.0
user=chromeuser
autorestart=true
priority=20
environment=HOME="/home/chromeuser",USER="chromeuser",LOGNAME="chromeuser"

[program:chrome]
command=/bin/bash -lc 'rm -f /data/chrome-profile/Singleton* /data/chrome-profile/DevToolsActivePort; exec /usr/bin/google-chrome --no-sandbox --test-type --disable-gpu --disable-software-rasterizer --disable-dev-shm-usage --remote-debugging-address=127.0.0.1 --remote-debugging-port=9223 --user-data-dir=/data/chrome-profile --disable-zygote --no-first-run --remote-allow-origins=* --window-size=1280,800 --window-position=0,0 about:blank'
directory=/home/chromeuser
environment=DISPLAY=":1",HOME="/home/chromeuser",USER="chromeuser",LOGNAME="chromeuser",XDG_CONFIG_HOME="/home/chromeuser/.config",XDG_CACHE_HOME="/home/chromeuser/.cache"
user=chromeuser
autorestart=true
priority=30
startsecs=5
stdout_logfile=/tmp/chrome_out.log
stderr_logfile=/tmp/chrome_err.log

[program:cdp_proxy]
command=/usr/bin/socat TCP-LISTEN:9222,fork,reuseaddr,bind=0.0.0.0 TCP:127.0.0.1:9223
autorestart=true
priority=40
stdout_logfile=/tmp/cdp_proxy_out.log
stderr_logfile=/tmp/cdp_proxy_err.log