# supervisord.conf
[unix_http_server]
file=/tmp/supervisor/supervisor.sock
chmod=0700
chown=chromeuser:chromeuser

[supervisord]
nodaemon=true
user=chromeuser
logfile=/tmp/supervisor/supervisord.log
pidfile=/tmp/supervisor/supervisord.pid
loglevel=info
silent=false
childlogdir=/tmp/supervisor
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor/supervisor.sock

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1280x800x24 -ac +extension GLX +render -noreset
user=chromeuser
autostart=true
autorestart=true
priority=10
stderr_logfile=/tmp/supervisor/xvfb_err.log
stdout_logfile=/tmp/supervisor/xvfb_out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
stopwaitsecs=10
environment=HOME="/home/chromeuser",USER="chromeuser",LOGNAME="chromeuser"

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -forever -shared -listen 0.0.0.0
user=chromeuser
autostart=true
autorestart=true
priority=20
stderr_logfile=/tmp/supervisor/x11vnc_err.log
stdout_logfile=/tmp/supervisor/x11vnc_out.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
stopasgroup=true
killasgroup=true
stopwaitsecs=10
environment=HOME="/home/chromeuser",USER="chromeuser",LOGNAME="chromeuser"

[program:chrome]
command=/bin/bash -lc 'rm -f /data/chrome-profile/Singleton* /data/chrome-profile/DevToolsActivePort; \
    exec /usr/bin/google-chrome \
    --disable-gpu \
    --remote-debugging-address=127.0.0.1 \
    --remote-debugging-port=9223 \
    --user-data-dir=/data/chrome-profile \
    --no-first-run \
    --remote-allow-origins=* \
    --window-size=1280,800 \
    --window-position=0,0 \
    --enable-features=NetworkService,NetworkServiceInProcess \
    --metrics-recording-only \
    --password-store=basic \
    --use-mock-keychain'
directory=/home/chromeuser
user=chromeuser
autostart=true
autorestart=true
priority=30
startsecs=5
stopwaitsecs=15
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/supervisor/chrome_out.log
stderr_logfile=/tmp/supervisor/chrome_err.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=5
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=5
environment=DISPLAY=":1",HOME="/home/chromeuser",USER="chromeuser",LOGNAME="chromeuser",XDG_CONFIG_HOME="/home/chromeuser/.config",XDG_CACHE_HOME="/home/chromeuser/.cache"

[program:cdp_proxy]
command=/usr/bin/socat TCP-LISTEN:9222,fork,reuseaddr,bind=0.0.0.0 TCP:127.0.0.1:9223
user=chromeuser
autostart=true
autorestart=true
priority=40
stopwaitsecs=10
stopasgroup=true
killasgroup=true
stdout_logfile=/tmp/supervisor/cdp_proxy_out.log
stderr_logfile=/tmp/supervisor/cdp_proxy_err.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=3
